# root/etc/services.d/vozila/run
#!/usr/bin/with-contenv sh

if [ -n "$SKIP_APP" ]; then
	exit 0
fi

exec s6-setuidgid $CONTAINER_USER node /app/main.js 
